# ansible/roles/base/files/postreboot-smoke.service
[Unit]
Description=Cash-Track post-reboot smoke check
After=docker.service network-online.target
Wants=network-online.target
ConditionPathExists=/usr/local/sbin/postreboot-smoke.sh

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/postreboot-smoke.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
